# .config/autostart/pcmanfm.desktop
[Desktop Entry]
Exec=pcmanfm -d --desktop
Hidden=false
Name=pcmanfm-desktop
Type=Application
